# Dockerfiles/syncsup/supervisor.d/discosrv.conf
[program:discosrv]
command=/usr/local/sbin/discosrv -key /home/%(ENV_APPUSER)s/key.pem -cert /home/%(ENV_APPUSER)s/cert.pem -listen %(ENV_LISTEN_DISCO)s
user=%(ENV_APPUSER)s
