# Dockerfiles/syncsup/supervisor.d/relaysrv.conf
[program:relaysrv]
command=/usr/local/sbin/relaysrv -keys /home/%(ENV_APPUSER)s -pools="" -listen %(ENV_LISTEN_RELAY)s
user=%(ENV_APPUSER)s
