# Dockerfiles/master/supervisor.d/apt-cacher.conf
[program:apt-cacher-ng]
command=/usr/sbin/apt-cacher-ng -c /etc/apt-cacher-ng ForeGround=1 CacheDir=/var/cache/apt-cacher-ng
user=apt-cacher-ng
group=apt-cacher-ng
autorestart=true
