# Dockerfiles/master/conf.d/consul-server.conf
[program:consul-server]
command = /usr/local/bin/consul agent -config-dir /etc/consul.d -data-dir /var/lib/consul
stopsignal = INT
