# Dockerfiles/postgres/supervisor.d/db_check_started.conf
[program:db_check_started]
command=/usr/local/bin/db_check_started.sh
user=root
autorestart=false
