# Dockerfiles/postgres/supervisor.d/dbcc.conf
[program:dbcc]
environment=APP_KEY=%(ENV_DBCC_KEY)s
command=/usr/local/bin/dbcc -p %(ENV_DBCC_PORT)s
user=postgres
autorestart=true
autostart=false
