# Dockerfiles/postgres/supervisor.d/dbcc.conf
[program:dbcc]
command=/usr/local/bin/dbcc --key %(ENV_DBCC_KEY)s 
user=postgres
autorestart=true
