# Dockerfiles/postgres/supervisor.d/postgres.conf
[program:postgres]
command=/usr/lib/postgresql/%(ENV_PG_MAJOR)s/bin/postgres -D %(ENV_PGDATA)s -c config_file=%(ENV_PGDATA)s/postgresql.conf
user=postgres
autorestart=true
stopsignal=TERM
