# Dockerfiles/pdns/supervisor.d/pdns.conf
[program:pdns]
command=/usr/sbin/pdns_server --master --daemon=no
autostart=false
