# Dockerfiles/pdns/supervisor.d/pdns.conf
[program:pdns]
command=/usr/sbin/pdns_server --master --daemon=no
autostart=false

[program:pdns-prep]
command=/usr/local/bin/dbinit_pdns.sh
autorestart=false
