# Dockerfiles/webhook/supervisor.d/webhook.conf
[program:webhook]
command=/usr/local/sbin/webhook -hooks /home/app/hooks.json -hotreload -ip 127.0.0.1 -nopanic -verbose
user=op
environment = HOME="/home/op", USER="op"
