# Dockerfiles/webhook/supervisor.d/webtail.conf
[program:webtail]
command=/usr/local/sbin/webtail --root /home/app/log --http_addr 0.0.0.0:8600 --host %(ENV_CONSUP_DOMAIN)s --lines 200
