# Dockerfiles/prosody/supervisor.d/prosody.conf
[program:prosody]
command=prosody
user=prosody