# Dockerfiles/pgrest/supervisor.d/dbinit.conf
[program:dbinit]
command=/usr/local/sbin/dbinit.sh
autorestart=false