# Dockerfiles/pgrest/supervisor.d/pgrest-prep.conf
[program:pgrest-prep]
command=/usr/local/sbin/postgrest_dbinit.sh
autorestart=false
